# pkg/common/mytest.mytest.timer
[Unit]
Description=Run MyTest periodically

[Timer]
OnUnitActiveSec=2s
OnBootSec=2s

[Install]
WantedBy=timers.target