# pkg/common/mytest.mytest.service
[Unit]
Description=My Test

[Service]
Type=oneshot
ExecStart=/usr/bin/mytest
